# beagle-os/overlay/etc/systemd/system/beagle-update-boot-scan.service
[Unit]
Description=Run a Beagle OS update scan on every boot
After=beagle-thin-client-prepare.service beagle-healthcheck.service network-online.target
Wants=beagle-thin-client-prepare.service beagle-healthcheck.service network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/beagle-update-client scan --force --auto-apply-if-idle

[Install]
WantedBy=multi-user.target
