# beagle-os/overlay/etc/systemd/system/beagle-endpoint-dispatch.timer
[Unit]
Description=Poll queued Beagle OS endpoint actions periodically

[Timer]
OnBootSec=20s
OnUnitActiveSec=15s
AccuracySec=5s
Unit=beagle-endpoint-dispatch.service

[Install]
WantedBy=timers.target
