# beagle-os/overlay/etc/systemd/system/beagle-endpoint-dispatch.service
[Unit]
Description=Dispatch queued Beagle OS endpoint actions
After=beagle-thin-client-prepare.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/beagle-endpoint-dispatch
