# beagle-os/overlay/etc/systemd/system/beagle-update-boot-scan.service
[Unit]
Description=Run a Beagle OS update scan on every boot
After=network-online.target pve-thin-client-prepare.service
Wants=network-online.target pve-thin-client-prepare.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/beagle-update-client scan --force --auto-apply-if-idle

[Install]
WantedBy=multi-user.target
