# beagle-os/overlay/etc/systemd/system/beagle-endpoint-report.service
[Unit]
Description=Report Beagle OS endpoint state to the control plane
After=pve-thin-client-prepare.service beagle-healthcheck.service network-online.target
Wants=pve-thin-client-prepare.service beagle-healthcheck.service network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/beagle-endpoint-report
