# thin-client-assistant/live-build/config/includes.chroot/etc/systemd/system/beagle-kiosk.target
[Unit]
Description=Beagle OS Gaming kiosk target
Wants=beagle-kiosk.service network-online.target
After=network-online.target
AllowIsolate=yes
