# server-installer/live-build/config/includes.chroot/etc/systemd/system/getty@tty1.service.d/zz-beagle-server-installer.conf
[Service]
Environment=LANG=en_US.UTF-8
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I linux
