# proxmox-host/systemd/pve-dcv-ui-reapply.path
[Unit]
Description=Watch Proxmox UI files and reapply PVE DCV Integration after updates

[Path]
PathChanged=/usr/share/pve-manager/index.html.tpl
PathChanged=/usr/share/pve-manager/js/pvemanagerlib.js
Unit=pve-dcv-ui-reapply.service

[Install]
WantedBy=multi-user.target
