# thin-client-assistant/live-build/config/includes.chroot/etc/systemd/system/ssh.service.d/10-beagle-hostkeys.conf
[Unit]
Wants=beagle-ssh-hostkeys.service
After=beagle-ssh-hostkeys.service
