# beagle-host/systemd/beagle-ui-reapply.service
[Unit]
Description=Reapply Beagle OS Proxmox UI hook
After=local-fs.target pveproxy.service
Wants=pveproxy.service
ConditionPathExists=/usr/share/pve-manager/index.html.tpl

[Service]
Type=oneshot
EnvironmentFile=-/etc/beagle/host.env
ExecStart=__INSTALL_DIR__/scripts/install-proxmox-ui-integration.sh
WorkingDirectory=__INSTALL_DIR__
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
