# beagle-host/systemd/beagle-control-plane.service
[Unit]
Description=Beagle OS Control Plane API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/beagle/beagle-manager.env
ExecStart=/usr/bin/python3 __INSTALL_DIR__/beagle-host/bin/beagle-control-plane.py
Restart=always
RestartSec=3
User=root
Group=root
UMask=0077
NoNewPrivileges=yes
PrivateDevices=no
PrivateTmp=yes
ProtectHome=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=no
# Proxmox provider needs kernel tunables (e.g. ip_forward); keep off
LockPersonality=yes
ProtectSystem=strict
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RuntimeDirectory=beagle-control-plane
RuntimeDirectoryMode=0750
ReadWritePaths=/var/lib/beagle /run/beagle-control-plane -/var/lib/vz -/etc/pve -/var/log/pve -/var/lib/libvirt/images -/var/lib/libvirt/qemu /etc/beagle/novnc
SystemCallArchitectures=native
SystemCallFilter=@system-service @network-io @file-system @basic-io
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_DAC_READ_SEARCH CAP_SETUID CAP_SETGID CAP_SYS_ADMIN
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
