# beagle-host/systemd/beagle-cluster-auto-join.service
[Unit]
Description=Beagle OS Cluster Auto-Join
After=network-online.target beagle-control-plane.service
Wants=network-online.target beagle-control-plane.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/beagle/host.env
EnvironmentFile=-/etc/beagle/beagle-manager.env
ExecStart=__INSTALL_DIR__/scripts/beagle-cluster-auto-join.sh
User=root
Group=root
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
ProtectSystem=strict
ReadWritePaths=/etc/beagle /var/lib/beagle

[Install]
WantedBy=multi-user.target
