# beagle-host/systemd/beagle-stream-reconciler.service
[Unit]
Description=Beagle Stream Reconciler (Python daemon)
After=network-online.target beagle-control-plane.service
Wants=network-online.target
PartOf=beagle-control-plane.service

[Service]
Type=simple
EnvironmentFile=-/etc/beagle/host.env
EnvironmentFile=-/etc/beagle/beagle-manager.env
ExecStart=/usr/bin/python3 __INSTALL_DIR__/beagle-host/services/stream_reconciler.py
WorkingDirectory=__INSTALL_DIR__
Restart=on-failure
RestartSec=30
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
# This service manages nftables and enables IPv4 forwarding.
ProtectKernelTunables=no
LockPersonality=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
