# beagle-host/systemd/beagle-novnc-proxy.service
[Unit]
Description=Beagle noVNC websockify proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=beagle-manager
Group=beagle-manager
ExecStart=/usr/bin/websockify --verbose --token-plugin beagle_novnc_token.BeagleTokenFile --token-source /etc/beagle/novnc/console-tokens.json 127.0.0.1:6080
Environment=PYTHONPATH=/opt/beagle/lib
Restart=always
RestartSec=2
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
ProtectSystem=strict
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
ReadWritePaths=/etc/beagle/novnc /opt/beagle/lib
SystemCallArchitectures=native
SystemCallFilter=@system-service @network-io @basic-io
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
