# beagle-host/systemd/beagle-novnc-proxy.service
[Unit]
Description=Beagle noVNC websockify proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/websockify --verbose --token-plugin TokenFile --token-source /etc/beagle/novnc/tokens 127.0.0.1:6080
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
