# beagle-host/systemd/beagle-wireguard-reconcile.service
[Unit]
Description=Reconcile Beagle WireGuard peers
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=__INSTALL_DIR__/scripts/apply-beagle-wireguard.sh --enable
User=root
Group=root
UMask=0077
NoNewPrivileges=yes
