# beagle-host/systemd/beagle-artifacts-watchdog.service
[Unit]
Description=Check Beagle hosted artifacts and trigger repair when needed
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=__INSTALL_DIR__/scripts/artifact-watchdog.sh
WorkingDirectory=__INSTALL_DIR__
User=root
Group=root
NoNewPrivileges=no
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ReadWritePaths=__INSTALL_DIR__/dist /var/lib/beagle /etc/beagle

[Install]
WantedBy=multi-user.target
