# beagle-host/systemd/beagle-artifacts-refresh.service
[Unit]
Description=Refresh Beagle OS hosted artifacts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/beagle/host.env
ExecStart=__INSTALL_DIR__/scripts/refresh-host-artifacts.sh
WorkingDirectory=__INSTALL_DIR__
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
