# beagle-os/overlay/etc/systemd/system/systemd-networkd-wait-online.service.d/beagle-any.conf
[Service]
# Beagle OS thin clients use wired ethernet (enp*). WiFi may not be present.
# --any: succeed as soon as at least one managed interface is online.
ExecStart=
ExecStart=/lib/systemd/systemd-networkd-wait-online --any
