# thin-client-assistant/live-build/config/includes.chroot/etc/systemd/system/getty@.service.d/zz-beagle-default.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- \u' --noclear - $TERM
