# server-installer/live-build/config/includes.chroot/etc/systemd/system/serial-getty@ttyS0.service.d/zz-beagle-server-installer.conf
[Service]
Environment=LANG=en_US.UTF-8
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear --keep-baud 115200,57600,38400,9600 %I vt220
