# beagle-kiosk/systemd/beagle-kiosk-update-catalog.service
# Beagle OS Gaming Kiosk - MIT Licensed
[Unit]
Description=Refresh Beagle OS Gaming catalog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/beagle-kiosk
ExecStart=/usr/bin/python3 /opt/beagle-kiosk/update_catalog.py
